# deploy/super_web.ini
[program:wanglibao_web]
command = /var/www/wanglibao/virt-wanglibao/bin/gunicorn -w 7 -k eventlet -b 0.0.0.0:8056 wanglibao.wsgi:application
directory = /var/www/wanglibao/wanglibao-backend
user = wangli
stdout_logfile = /tmp/web_super.log
stderr_logfile = /tmp/web_super_err.log
