# deploy/super_task.ini

[program:wanglibao_task]
command = /var/www/wanglibao/virt-wanglibao/bin/python manage.py celeryd -l info -E -B
directory = /var/www/wanglibao/wanglibao-backend
user = wangli
stdout_logfile = /tmp/task_super.log
stderr_logfile = /tmp/task_super_err.log

[program:wanglibao_task_cam]
command = /var/www/wanglibao/virt-wanglibao/bin/python manage.py celerycam
directory = /var/www/wanglibao/wanglibao-backend
user = wangli
stdout_logfile = /tmp/cam_super.log
stderr_logfile = /tmp/cam_super_err.log
